# deploy/reactor_ctrl_db_backup.service
[Unit]
Description=reactor_ctrl database backup (one-shot SQL dump)
After=network-online.target mariadb.service mysql.service
Wants=network-online.target

[Service]
Type=oneshot
User=pthuerlemann
Group=pthuerlemann
WorkingDirectory=/home/pthuerlemann/reactor_ctrl
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/pthuerlemann/reactor_ctrl/.venv/bin/python /home/pthuerlemann/reactor_ctrl/run_database_backup.py
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
